# launchd/com.claude.telegram-bot.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.claude.telegram-bot</string>

    <key>ProgramArguments</key>
    <array>
        <!-- UPDATE: Path to your Python -->
        <string>/usr/bin/python3</string>
        <!-- UPDATE: Path to telegram-bot.py -->
        <string>/Users/YOUR_USERNAME/claude-code-telegram/telegram-bot.py</string>
    </array>

    <key>WorkingDirectory</key>
    <!-- UPDATE: Path to repo -->
    <string>/Users/YOUR_USERNAME/claude-code-telegram</string>

    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/usr/local/bin:/usr/bin:/bin:/opt/homebrew/bin</string>
    </dict>

    <key>RunAtLoad</key>
    <true/>

    <key>KeepAlive</key>
    <true/>

    <key>StandardOutPath</key>
    <string>/tmp/telegram-bot.log</string>

    <key>StandardErrorPath</key>
    <string>/tmp/telegram-bot.error.log</string>
</dict>
</plist>
